# tomcat/tomcat7/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/bin/bash -c "exec ${CATALINA_HOME}/bin/catalina.sh run"